# deploy/presage-serve.service
[Unit]
Description=PRESAGE dashboard (read-only, boucle de render)
After=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/mes-bots-finance
ExecStart=%h/mes-bots-finance/venv/bin/python -m dashboard.serve
Restart=always
RestartSec=10
Environment=PRESAGE_PORT=8000

[Install]
WantedBy=default.target
