# deploy/launchd_mac/com.olivier.presage-sync-from-vm.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.olivier.presage-sync-from-vm</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/olivierlegendre/mes-bots-finance/scripts/sync_db_from_hetzner.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/olivierlegendre/mes-bots-finance</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Minute</key>
        <integer>5</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/olivierlegendre/mes-bots-finance/logs/launchd_sync_from_vm.out</string>

    <key>StandardErrorPath</key>
    <string>/Users/olivierlegendre/mes-bots-finance/logs/launchd_sync_from_vm.err</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
